# config/crontab
# claude-channels scheduled tasks — fires webhooks into localhost:8788

# User's Daily briefing at 7:30am
30 7 * * * curl -s -X POST http://127.0.0.1:8788/briefing -H 'Content-Type: application/json' -d '{"task":"daily_briefing"}'

# User's Check-in every hour (6am-9pm Pacific)
0 6-21 * * * curl -s -X POST http://127.0.0.1:8788/reconcile -H 'Content-Type: application/json' -d '{"task":"reconcile"}'

# User's Weather check at 7am and 5pm
0 7,17 * * * curl -s -X POST http://127.0.0.1:8788/weather -H 'Content-Type: application/json' -d '{"task":"weather"}'

# User's End-of-day wrap-up at 11pm
0 23 * * * curl -s -X POST http://127.0.0.1:8788/eod -H 'Content-Type: application/json' -d '{"task":"end_of_day"}'
